# com.example.housing-monitor-healthcheck.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.YOURUSER.housing-monitor-healthcheck</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/path/to/housing-monitor/healthcheck.py</string>
    </array>
    <key>StartInterval</key>
    <integer>345600</integer>
    <key>StandardOutPath</key>
    <string>/path/to/housing-monitor/monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/housing-monitor/monitor.log</string>
    <key>WorkingDirectory</key>
    <string>/path/to/housing-monitor</string>
</dict>
</plist>
